# swimmy.service
[Unit]
Description=Swimmy Trading System
Documentation=https://github.com/swimmy/swimmy
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=%h/swimmy/run.sh
Restart=always
RestartSec=10

# Environment
Environment=HOME=%h
EnvironmentFile=%h/swimmy/.env

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h/swimmy/doc/logs
ReadWritePaths=%h/swimmy/.opus
ReadWritePaths=%h/swimmy/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=swimmy

[Install]
WantedBy=default.target
