# projects/auto_journal/auto-journal.service
[Unit]
Description=Swimmy Auto-Journal Daemon
After=network.target

[Service]
Type=simple
User=swimmy
WorkingDirectory=%h/swimmy
Environment=DISPLAY=:0
ExecStart=%h/swimmy/.venv/bin/python -m auto_journal.auto_journal daemon
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
