# systemd/data_keeper.service
[Unit]
Description=Swimmy Data Keeper Service (Article 5 Compliant)
After=network.target

[Service]
Type=simple
User=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/data_keeper.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/swimmy/swimmy/logs/data_keeper.log
StandardError=append:/home/swimmy/swimmy/logs/data_keeper.log

[Install]
WantedBy=default.target
