# systemd/mt5_account_sync.service
[Unit]
Description=Swimmy MT5 Account Sync Service
Documentation=https://github.com/gethinu/swimmy-fx-tribe
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/usr/bin/python3 /home/swimmy/swimmy/tools/mt5_account_sync.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment
EnvironmentFile=/home/swimmy/swimmy/.env

# Logging
StandardOutput=append:/home/swimmy/swimmy/logs/mt5_account_sync.log
StandardError=append:/home/swimmy/swimmy/logs/mt5_account_sync.log

# Failure notification hook (sends to Discord via wrapper script)
ExecStopPost=/bin/bash /home/swimmy/swimmy/tools/notify_failure.sh $$SERVICE_RESULT "MT5 Account Sync"

[Install]
WantedBy=default.target
