# systemd/pending_manager.service
[Unit]
Description=Swimmy Pending Manager (López de Prado)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=/usr/bin/python3 %h/swimmy/tools/pending_manager.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:%h/swimmy/logs/pending_manager.log
StandardError=append:%h/swimmy/logs/pending_manager.log

[Install]
WantedBy=default.target
