# systemd/strategy_hunter.service
[Unit]
Description=Swimmy Strategy Hunter
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/usr/bin/python3 /home/swimmy/swimmy/tools/strategy_hunter.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/swimmy/swimmy/logs/strategy_hunter.log
StandardError=append:/home/swimmy/swimmy/logs/strategy_hunter.log

[Install]
WantedBy=default.target
