# systemd/strategy_hunter.service
[Unit]
Description=Swimmy Strategy Hunter
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=/usr/bin/python3 %h/swimmy/tools/strategy_hunter.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:%h/swimmy/logs/strategy_hunter.log
StandardError=append:%h/swimmy/logs/strategy_hunter.log

[Install]
WantedBy=default.target
