# systemd/swimmy-data-keeper.service
[Unit]
Description=Swimmy Data Keeper Service (Article 5 Compliant)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/data_keeper.py
Restart=always
RestartSec=10
# P4 (Thread B): cap the data-keeper's resident memory via the cgroup so a
# runaway read/merge cannot starve the brain (3072MB) and school (4096MB)
# heaps on the shared ops box. On OOM the cgroup kills only this unit, which
# Restart=always then relaunches cleanly.
MemoryMax=3G
Environment=HOME=/home/swimmy
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/swimmy/swimmy/.env
StandardOutput=append:/home/swimmy/swimmy/logs/data_keeper.log
StandardError=append:/home/swimmy/swimmy/logs/data_keeper.log

[Install]
WantedBy=default.target
