# systemd/swimmy-brain.service
[Unit]
Description=Swimmy Brain Service
Documentation=https://github.com/swimmy/swimmy
After=network.target

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/run.sh
Restart=always
RestartSec=10

# Environment
Environment=HOME=/home/swimmy
Environment=SWIMMY_SBCL_DYNAMIC_SPACE_MB=4096
EnvironmentFile=/home/swimmy/swimmy/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=swimmy-brain

[Install]
WantedBy=default.target
