# systemd/swimmy-risk.service
[Unit]
Description=Swimmy Risk Gateway Service (Graham)
After=network.target

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/risk_gateway.py
Restart=always
RestartSec=5
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env
# Logs
StandardOutput=append:/home/swimmy/swimmy/logs/risk.log
StandardError=append:/home/swimmy/swimmy/logs/risk.log

[Install]
WantedBy=default.target
