# systemd/swimmy-oos-monitor.timer
[Unit]
Description=Run Swimmy OOS Monitor every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
