# systemd/swimmy-school.service
[Unit]
Description=Swimmy School Service
Documentation=https://github.com/swimmy/swimmy
After=network.target

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStartPre=/usr/bin/sbcl --disable-debugger --script /home/swimmy/swimmy/tools/restore_legend_61.lisp
ExecStart=/usr/bin/sbcl --dynamic-space-size ${SWIMMY_SBCL_DYNAMIC_SPACE_MB} --noinform --load /home/swimmy/swimmy/src/lisp/school/school-daemon.lisp
Restart=always
RestartSec=10

# Environment
Environment=HOME=/home/swimmy
Environment=SWIMMY_SBCL_DYNAMIC_SPACE_MB=4096
EnvironmentFile=/home/swimmy/swimmy/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=swimmy-school

[Install]
WantedBy=default.target
