# systemd/swimmy-openclaw-signal-sync.service
[Unit]
Description=Swimmy OpenClaw Signal Sync
After=network.target

[Service]
Type=oneshot
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/openclaw_signal_sync.py
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env
Environment=POLYCLAW_USE_HEURISTIC_IF_NO_OPENCLAW_CMD=1
Environment=POLYCLAW_HEURISTIC_QUESTION_KEYWORDS=nba,nfl,mlb,nhl,ncaa,lakers,warriors,chiefs,yankees
Environment=POLYCLAW_SIGNALS_FILE=/home/swimmy/swimmy/data/openclaw/signals.jsonl
Environment=POLYCLAW_SIGNALS_META_FILE=/home/swimmy/swimmy/data/openclaw/signals.meta.json
Environment=POLYCLAW_SIGNAL_SYNC_MIN_SIGNALS=1
Environment=POLYCLAW_SIGNAL_SYNC_TIMEOUT_SECONDS=30

StandardOutput=append:/home/swimmy/swimmy/logs/openclaw_signal_sync.log
StandardError=append:/home/swimmy/swimmy/logs/openclaw_signal_sync.log
