# systemd/swimmy-pattern-backend-calibration.timer
[Unit]
Description=Run Swimmy Pattern backend calibration daily

[Timer]
OnCalendar=*-*-* 06:20:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
