# systemd/swimmy-discord-bot.service
[Unit]
Description=Swimmy Discord Bot Service
After=network.target

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/src/python/discord_bot.py
Restart=always
RestartSec=5
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env
# Logs
StandardOutput=append:/home/swimmy/swimmy/logs/discord_bot.log
StandardError=append:/home/swimmy/swimmy/logs/discord_bot.log

[Install]
WantedBy=default.target
