# systemd/swimmy-guardian.service
[Unit]
Description=Swimmy Guardian Service (Watchdog)
Documentation=https://github.com/swimmy/swimmy
After=network.target

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/watchdog.py
Restart=always
RestartSec=10

# Environment
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=swimmy-guardian

[Install]
WantedBy=default.target
