# systemd/swimmy-backtest.service
[Unit]
Description=Swimmy Backtest Service (Naval Modularization)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/backtest_service.py
Restart=on-failure
RestartSec=10
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env
LimitNOFILE=65535
# Logs
StandardOutput=append:/home/swimmy/swimmy/logs/backtest.log
StandardError=append:/home/swimmy/swimmy/logs/backtest.log

[Install]
WantedBy=default.target
