# systemd/swimmy-polymarket-openclaw-status.timer
[Unit]
Description=Run Swimmy Polymarket OpenClaw status monitor every 10 minutes

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
