# systemd/swimmy-history-update.timer
[Unit]
Description=Run Swimmy Historical Data Update Monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
