# systemd/swimmy-evolution.service
[Unit]
Description=Swimmy Evolution Daemon (Hyper-Time)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/evolution_daemon.py
Restart=always
RestartSec=10
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env
Environment=PYTHONUNBUFFERED=1
Nice=19
CPUSchedulingPolicy=batch
StandardOutput=append:/home/swimmy/swimmy/logs/evolution_daemon.log
StandardError=append:/home/swimmy/swimmy/logs/evolution_daemon.log

[Install]
WantedBy=default.target
