# systemd/swimmy-system-audit.timer
[Unit]
Description=Run Swimmy System Audit daily

[Timer]
OnCalendar=*-*-* 06:40:00
Persistent=true
RandomizedDelaySec=20min

[Install]
WantedBy=timers.target
