# systemd/swimmy-pattern-similarity.service
[Unit]
Description=Swimmy Pattern Similarity Service (REQ/REP 5564, S-exp)
After=network.target

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/pattern_similarity_service.py
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
Environment=HOME=/home/swimmy
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/swimmy/swimmy/.env
StandardOutput=append:/home/swimmy/swimmy/logs/pattern_similarity.log
StandardError=append:/home/swimmy/swimmy/logs/pattern_similarity.log

[Install]
WantedBy=default.target

