# systemd/swimmy-archive-reconcile.timer
[Unit]
Description=Run Swimmy Archive Drift Reconcile hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target

