# systemd/swimmy-history-update.service
[Unit]
Description=Swimmy Historical Data Update (Request to Guardian)
After=network.target swimmy-guardian.service

[Service]
Type=oneshot
User=swimmy
Group=swimmy
PermissionsStartOnly=true
WorkingDirectory=/home/swimmy/swimmy
ExecStartPre=-/usr/bin/systemctl stop swimmy-guardian.service
ExecStart=/home/swimmy/swimmy/.venv/bin/python3 /home/swimmy/swimmy/tools/update_history_smart.py
ExecStartPost=/usr/bin/systemctl start swimmy-guardian.service
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env
StandardOutput=append:/home/swimmy/swimmy/logs/history_update.log
StandardError=append:/home/swimmy/swimmy/logs/history_update.log
