# systemd/swimmy-evolution.service
[Unit]
Description=Swimmy Evolution Daemon (Hyper-Time)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=/usr/bin/python3 %h/swimmy/tools/evolution_daemon.py
Restart=always
RestartSec=10
EnvironmentFile=%h/swimmy/.env
Environment=PYTHONUNBUFFERED=1
Nice=19
CPUSchedulingPolicy=batch
StandardOutput=append:%h/swimmy/logs/evolution_daemon.log
StandardError=append:%h/swimmy/logs/evolution_daemon.log

[Install]
WantedBy=default.target
