# systemd/swimmy-guardian.service
[Unit]
Description=Swimmy Guardian Service (Rust)
Documentation=https://github.com/swimmy/swimmy
After=network.target

[Service]
Type=simple
User=swimmy
Group=swimmy
WorkingDirectory=/home/swimmy/swimmy
# Phase 1: (disabled) pkill was killing backtest-only guardian and causing restart loops
# Phase 2: Wait for OS to release sockets
ExecStartPre=/bin/sleep 2
# Phase 3: Force-release ZMQ ports (5557=Guardian PUB)
ExecStartPre=-/usr/bin/fuser -k 5557/tcp
# Phase 4: Final cooldown
ExecStartPre=/bin/sleep 1
ExecStart=/home/swimmy/swimmy/guardian/target/release/guardian
StandardOutput=append:/home/swimmy/swimmy/logs/guardian.log
StandardError=append:/home/swimmy/swimmy/logs/guardian.log
Environment=HOME=/home/swimmy
EnvironmentFile=/home/swimmy/swimmy/.env
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=default.target
