# tools/systemd/swimmy-keeper.service
[Unit]
Description=Swimmy Data Keeper Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=%h/swimmy/.venv/bin/python3 tools/data_keeper.py
Restart=always
RestartSec=5
EnvironmentFile=%h/swimmy/.env
# Logs
StandardOutput=append:%h/swimmy/logs/keeper.log
StandardError=append:%h/swimmy/logs/keeper.log

[Install]
WantedBy=default.target
