# tools/systemd/swimmy-backtest.service
[Unit]
Description=Swimmy Backtest Service (Naval Modularization)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=%h/swimmy/.venv/bin/python3 tools/backtest_service.py
Restart=always
RestartSec=5
EnvironmentFile=%h/swimmy/.env
LimitNOFILE=65535
# Logs
StandardOutput=append:%h/swimmy/logs/backtest.log
StandardError=append:%h/swimmy/logs/backtest.log

[Install]
WantedBy=default.target
