# tools/systemd/swimmy-notifier.service
[Unit]
Description=Swimmy Discord Notifier Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=%h/swimmy/.venv/bin/python3 tools/notifier.py
Restart=always
RestartSec=5
EnvironmentFile=%h/swimmy/.env
# Logs
StandardOutput=append:%h/swimmy/logs/notifier.log
StandardError=append:%h/swimmy/logs/notifier.log

[Install]
WantedBy=default.target
