# tools/systemd/swimmy-risk.service
[Unit]
Description=Swimmy Risk Gateway Service (Graham)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/swimmy
ExecStart=%h/swimmy/.venv/bin/python3 tools/risk_gateway.py
Restart=always
RestartSec=5
EnvironmentFile=%h/swimmy/.env
# Logs
StandardOutput=append:%h/swimmy/logs/risk.log
StandardError=append:%h/swimmy/logs/risk.log

[Install]
WantedBy=default.target
